# HOME/config/systemd/user/alarms.timer
[Unit]
Description=Run alarm every 3min

[Timer]
OnBootSec=10min
OnUnitActiveSec=3min

[Install]
WantedBy=timers.target
