# HOME/config/systemd/user/alarms.service
[Unit]
Description=Alarms

[Service]
ExecStart=%h/.local/bin/alarms.py

[Install]
WantedBy=default.target
