# HOME/config/systemd/user/poweroff.service
[Unit]
Description=PowerOff system

[Service]
ExecStart=systemctl poweroff

[Install]
WantedBy=default.target
