# HOME/config/systemd/user/poweroff.timer
[Unit]
Description=Poweroff system

[Timer]
OnCalendar=*-*-* 22:00:00 CST

[Install]
WantedBy=timers.target
